# systemd/reminders.service
[Unit]
Description=Periodic reminders
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bob
WorkingDirectory=/home/bob/clock
ExecStart=/home/bob/go/bin/reminders
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target