# systemd/weather.service
[Unit]
Description=Weather notification
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bob
WorkingDirectory=/home/bob/clock
ExecStart=/home/bob/go/bin/weather
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target